# ubuntu/usensord.conf
description "usensord service for haptics support"

start on started dbus
stop on stopped dbus

env VIBRATE_SCALE=15

exec usensord

